# test_server/etc/supervisor/evs_test_server.conf
[program:test_eventsource]
command=/usr/local/www/testevs.amvtek.com/bin/twistd 
    -n --pidfile= 
    test_eventsource
    --host=127.0.0.1
    --port=7676
user=eventsource
directory=/usr/local/www/testevs.amvtek.com/EventSource/test_server
process_name=%(program_name)s
num_procs=1
